# gateway/systemd/litellm.service
# ============================================================
# FOTA 智能诊断平台 - LiteLLM Gateway systemd 服务配置
# ============================================================
# 安装步骤：
# 1. sudo cp systemd/litellm.service /etc/systemd/system/
# 2. sudo systemctl daemon-reload
# 3. sudo systemctl enable litellm
# 4. sudo systemctl start litellm

[Unit]
Description=LiteLLM Proxy - FOTA LLM API Gateway
Documentation=https://docs.litellm.ai/docs/proxy/deploy
After=network.target

[Service]
Type=exec
User=litellm
Group=litellm
WorkingDirectory=/opt/litellm-proxy

# 加载环境变量
EnvironmentFile=/opt/litellm-proxy/.env

# 启动命令
ExecStart=/opt/litellm-proxy/venv/bin/litellm \
    --config /opt/litellm-proxy/config.yaml \
    --host 127.0.0.1 \
    --port 4000 \
    --num_workers 4

# 日志输出到 journald（可用 journalctl -u litellm 查看）
StandardOutput=journal
StandardError=journal

# 自动重启
Restart=always
RestartSec=5

# 安全加固
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/litellm-proxy/logs
PrivateTmp=yes

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
